# plugin/source/dynamix.unraid.net/etc/cron.daily/MyServersGetJS.sh
#!/bin/bash
# Check daily for updated My Servers JavaScript web components e.g. the user profile component
bash /usr/local/emhttp/plugins/dynamix.my.servers/scripts/getJS.sh &> /dev/null
