# plugin/source/dynamix.unraid.net/etc/profile.d/bun.sh
#! /bin/bash

# Add alias for bunx
alias bunx='bun x'
